# dot_config/systemd/systemd/user/hipkid.service
[Unit]
Description=HiPKI Local Sign Server
After=default.target

[Service]
ExecStart=%h/.local/bin/hipki.sh

[Install]
WantedBy=default.target


